# custom/stable-diffusion.service
[Unit]
Description=Stable-Diffusion Script

[Service]
ExecStart=/workspace/stable-diffusion-webui/webui.sh

[Install]
WantedBy=multi-user.target